# scripts/systemd/research-kb-sync.service
[Unit]
Description=Research-KB KuzuDB Sync
Documentation=https://github.com/brandon-behring/research-kb
After=postgresql.service research-kb-daemon.service

[Service]
Type=oneshot
ExecStart=%h/Claude/research-kb/venv/bin/python %h/Claude/research-kb/scripts/sync_kuzu.py
WorkingDirectory=%h/Claude/research-kb
Environment=PATH=%h/Claude/research-kb/venv/bin:/usr/local/bin:/usr/bin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=research-kb-sync

# Safety
TimeoutStartSec=300
Nice=10

[Install]
WantedBy=default.target
