# services/research-kb-daemon.service
[Unit]
Description=Research KB Daemon - Low-latency query service
Documentation=https://github.com/brandon-behring/research-kb
After=postgresql.service
# Wants embed server to be running (soft dependency)
Wants=research-kb-embed.service

[Service]
Type=simple
ExecStart=<PROJECT_ROOT>/venv/bin/research-kb-daemon
Restart=on-failure
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=EMBED_SOCKET_PATH=/tmp/research_kb_embed.sock

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

# Allow write to /tmp for socket
ReadWritePaths=/tmp

[Install]
WantedBy=default.target
